# omarchy-zed-sync.service
[Unit]
Description=Omarchy Zed Theme Sync
Documentation=https://github.com/aps6/omarchy-theme-zed
After=graphical-session.target xdg-desktop-portal.service
Wants=graphical-session.target
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=simple
ExecStart=%h/.local/bin/omarchy-zed-theme-sync.sh watch
Restart=on-failure
RestartSec=10
Environment=HOME=%h
RestartPreventExitStatus=1

[Install]
WantedBy=graphical-session.target
